# ops/launchd/com.symphony.bob-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.symphony.bob-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bob/AI-Server/scripts/bob-watchdog.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>BOB_REPO_DIR</key><string>/Users/bob/AI-Server</string>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/bob/AI-Server/data/task_runner/bob-watchdog.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bob/AI-Server/data/task_runner/bob-watchdog.err.log</string>
</dict>
</plist>
