# scripts/launchd/com.symphony.audio-intake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.symphony.audio-intake</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/bob/AI-Server/scripts/audio_intake_worker.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>AI_SERVER_ROOT</key>
    <string>/Users/bob/AI-Server</string>
    <key>CORTEX_URL</key>
    <string>http://127.0.0.1:8102</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/bob/AI-Server/data/audio_intake/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bob/AI-Server/data/audio_intake/launchd.err.log</string>
</dict>
</plist>
