# setup/launchd/com.symphony.trading-topic-graph-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.trading-topic-graph-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bob/AI-Server/trading/cross_topic_graph_runner.py</string>
        <string>--days</string>
        <string>14</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bob/AI-Server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bob/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bob</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bob/AI-Server/logs/trading-topic-graph-daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bob/AI-Server/logs/trading-topic-graph-daily.err</string>
</dict>
</plist>
