# setup/launchd/com.symphonysh.dropbox-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphonysh.dropbox-organizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/bob/AI-Server/scripts/dropbox-organizer.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/dropbox-organizer.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/dropbox-organizer.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DROPBOX_ROOT</key>
        <string>/Users/bob/Dropbox</string>
    </dict>
</dict>
</plist>
