# setup/launchd/com.symphony.learning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.learning</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/bob/AI-Server/openclaw/continuous_learning.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DATA_DIR</key>
        <string>/Users/bob/AI-Server/data</string>
        <key>SYMPHONY_ROOT</key>
        <string>/Users/bob/AI-Server</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/symphony-learning.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/symphony-learning.log</string>
</dict>
</plist>
