# setup/launchd/com.symphony.bob-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.bob-maintenance</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/bob/AI-Server/.venv/bin/python3</string>
        <string>/Users/bob/AI-Server/tools/bob_maintenance.py</string>
        <string>--run</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/bob/AI-Server</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/bob/AI-Server/logs/bob-maintenance.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bob/AI-Server/logs/bob-maintenance.err</string>
</dict>
</plist>
