# setup/launchd/com.symphony.graph-drift-watcher-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.graph-drift-watcher-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bob/AI-Server/trading/graph_drift_watcher.py</string>
        <string>--edge-delta</string>
        <string>0.20</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/bob/AI-Server</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bob/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bob</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/bob/AI-Server/logs/graph-drift-watcher-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bob/AI-Server/logs/graph-drift-watcher-daily.err</string>
</dict>
</plist>
