# setup/launchd/com.symphony.imessage-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.symphony.imessage-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/bob/AI-Server/.venv-imessage/bin/python3</string>
    <string>/Users/bob/AI-Server/scripts/imessage-server.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/bob/AI-Server</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/imessage-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/imessage-bridge.log</string>
</dict>
</plist>
