# setup/launchd/com.symphony.focus-ops-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.focus-ops-monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bob/AI-Server/tools/focus_ops_monitor.py</string>
        <string>--once</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/bob/AI-Server</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bob/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bob</string>
    </dict>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/bob/AI-Server/logs/focus-ops-monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bob/AI-Server/logs/focus-ops-monitor.err</string>
</dict>
</plist>
