# setup/launchd/com.symphony.trading-provider-slo-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.trading-provider-slo-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bob/AI-Server/trading/provider_slo_monitor.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bob/AI-Server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bob/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/bob</string>
        <key>TRADING_SLO_WINDOW_MINUTES</key>
        <string>60</string>
        <key>TRADING_SLO_MIN_SAMPLES</key>
        <string>8</string>
        <key>TRADING_SLO_FALLBACK_RATE</key>
        <string>0.35</string>
        <key>TRADING_SLO_P95_LATENCY_MS</key>
        <string>9000</string>
    </dict>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bob/AI-Server/logs/trading-provider-slo-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bob/AI-Server/logs/trading-provider-slo-monitor.err</string>
</dict>
</plist>
