# setup/launchd/com.symphony.overnight-learner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.symphony.overnight-learner</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bob/AI-Server/tools/overnight_learner.py</string>
        <string>--all</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/bob/AI-Server</string>
    
    <!-- Run at 11 PM every night -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/bob/AI-Server/logs/overnight-learner.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/bob/AI-Server/logs/overnight-learner.err</string>
</dict>
</plist>
