# scripts/com.symphony.imessage-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.symphony.imessage-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>REPLACE_WITH_VENV_PYTHON</string>
    <string>REPLACE_WITH_IMESSAGE_SERVER_PY</string>
  </array>
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_REPO_ROOT</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/imessage-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/imessage-bridge.log</string>
</dict>
</plist>
