# scripts/com.symphony.imessage-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  1. Adjust paths if your home directory is not /Users/bob
  2. cp scripts/com.symphony.imessage-bridge.plist.example ~/Library/LaunchAgents/com.symphony.imessage-bridge.plist
  3. launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.symphony.imessage-bridge.plist
  4. launchctl kickstart -k gui/$(id -u)/com.symphony.imessage-bridge
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.symphony.imessage-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/bob/AI-Server/.venv-imessage/bin/python3</string>
    <string>/Users/bob/AI-Server/scripts/imessage-server.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/bob/AI-Server</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/imessage-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/imessage-bridge.log</string>
</dict>
</plist>
